# src/files/etc/init.d/outlineGate
#!/bin/sh /etc/rc.common
# Example script
# Copyright (C) 2007 OpenWrt.org

USE_PROCD=1

START=96
STOP=96
LOGFILE="/var/log/outline-gate.log"
# PROCESS_LOG="/var/log/chisel_process.log"
CHISEL_BIN="/usr/bin/chisel"
[ -x "$CHISEL_BIN" ] || CHISEL_BIN="/tmp/chisel"

#init_process_log() {
#    [ -f "$PROCESS_LOG" ] || { touch "$PROCESS_LOG" && chmod 600 "$PROCESS_LOG"; }
#    if [ $(stat -c%s "$PROCESS_LOG" 2>/dev/null || echo 0) -gt 1000000 ]; then
#        mv "$PROCESS_LOG" "$PROCESS_LOG.1" 2>/dev/null
#        touch "$PROCESS_LOG" && chmod 600 "$PROCESS_LOG"
#    fi
#}

#log_msg() {
#    init_process_log
#    echo "$(date '+%Y-%m-%d %H:%M:%S') $1" >> "$PROCESS_LOG"
#}


start_service() {

        # log_msg "Starting outlineGate service"
        key=$(uci get routro.ireach.user)
        port=$(uci get routro.ireach.port)
        host=$(uci get routro.ireach.host)
        password=$(uci get routro.ireach.pass)

        mport=$(uci get routro.ireach.outlineport)

        oport=$(uci get routro.outlinegate.originport)
        ohost=$(uci get routro.outlinegate.originhost)


        if [ -z "$key" ] || [ -z "$port" ] || [ -z "$host" ] || [ -z "$password" ]  || [ -z "$mport" ] || [ -z "$oport" ] || [ -z "$ohost" ];then
                echo "outlineGate_status=MISSINF_CONFIG" > /tmp/outlineGate_status
                echo "Connect - service is not ready" >> $LOGFILE
                # log_msg "Missing configuration, aborting start"
                exit 0;
        fi

	procd_open_instance
        procd_set_param command /bin/sh -c "$CHISEL_BIN client --auth $key:$password ${host}:${port} R:0.0.0.0:$mport:$ohost:$oport >> $LOGFILE 2>&1"
    
         # Enable respawn if needed
        procd_set_param respawn

        procd_set_param stdout 1
        procd_set_param stderr 1
        procd_close_instance
        # log_msg "outlineGate service started"
}

stop() {

        echo "Connect - Stopping outlinegate service" >> $LOGFILE
        # log_msg "Stopping outlineGate service"
        # commands to kill application
}
